# cliptube.service
[Unit]
Description=cliptube directory watcher

[Service]
Type=forking
ExecStart=/home/chris/.local/bin/ddirwatch

[Install]
WantedBy=default.target
