# python-cliptube.service
[Unit]
Description=Clipboard watcher

[Install]
WantedBy=default.target

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/gcbwatch
