# configs/dirwatch.service
[Unit]
Description=cliptube youtube directory watcher

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/chris/.local/bin/dirwatch

[Install]
WantedBy=default.target

